# infra/systemd/mistral-mcp.service
[Unit]
Description=Mistral PVE MCP HTTP server
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/mistral/secrets.env
ExecStart=/usr/local/bin/mistral mcp --http
Restart=on-failure
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
