# infra/systemd/mistral-web.service
[Unit]
Description=Mistral PVE web settings UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/mistral/secrets.env
Environment=MISTRAL_CONFIG=/etc/mistral/config.yaml
Environment=MISTRAL_STATE=/var/lib/mistral/state.json
ExecStart=/usr/local/bin/mistral web
Restart=on-failure
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
